# files/speedtest_cron.sh
#!/bin/sh
. /etc/rc.subr
configctl speedtest run
